# crontab
MAILTO=andres@krutt.org
@hourly cat ~/.keyring/andreser@mit.edu.txt | kinit -l7d -f andreser
