# systemd/user/desktop.target
[Unit]
Description=My desktop system
Requires=default.target
After=default.target
AllowIsolate=yes
