# systemd/user/pyhesiodfs.service
[Service]
ExecStart=/usr/bin/pyhesiodfs -f /mit -o nonempty
 
[Install]
WantedBy=default.target

